# chroot/etc/udev/rules.d/84-ebs.rules
KERNEL=="sd*[!0-9]", ACTION=="add|change", RUN+="/home/marat/workspace/scalarizr-trunk/src/scalarizr/core/udev.py"